# deploy/systemd/uaft-web.service
[Unit]
Description=UAFT Web Application (Next.js)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/uaft
Environment="NODE_ENV=production"
Environment="PORT=3000"
EnvironmentFile=/var/www/uaft/.env.production
ExecStart=/usr/bin/node /var/www/uaft/.next/standalone/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=uaft-web

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/uaft/uploads /var/www/uaft/.next/cache

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
